# llm-stack/crontab/crontab
# 改行コードは"LF"にすること
* * * * * /bin/bash -l -c 'echo `date` >> /var/log/cron.log' 2>&1
* * * * * /usr/local/bin/python3 /root/main.py >> /var/log/cron.log 2>&1
